# supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n

[program:postgres]
command=/docker-entrypoint.sh postgres
